# lumen/config/supervisor/supervisor_ng_ams_pre_prod.conf
# TODO: Rename the following worker to be consistent with Bundle Management Gateway naming convention, as well the log filename.
# https://jira.mgcorp.co/browse/BG-34654
[program:purchase-gateway-pre-prod-worker-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /home/app_billing/applications/purchase-gateway/pre-prod/current/artisan queue:work
autostart=true
autorestart=true

numprocs=1
redirect_stderr=true
stdout_logfile=/home/app_billing/applications/purchase-gateway/pre-prod/logs/purchase-gateway-worker-postback-queue.log

[program:purchase-gateway-pre-prod-worker-email-send]
process_name=%(program_name)s_%(process_num)02d
command=php /home/app_billing/applications/purchase-gateway/pre-prod/current/artisan ng:worker send-email run
autostart=true
autorestart=true

numprocs=1
# startsecs
# This is the minimum number of seconds that supervisor expects the process to run for it to be considered successful
# If the process stops in under the amount of seconds specified, supervisor will consider it exited with error and
# it will throw a FATAL.
startsecs=1
redirect_stderr=true
stdout_logfile=/home/app_billing/applications/purchase-gateway/pre-prod/logs/purchase-gateway-worker.log

[program:purchase-gateway-pre-prod-worker-bundle-addons]
process_name=%(program_name)s_%(process_num)02d
command=php /home/app_billing/applications/purchase-gateway/pre-prod/current/artisan  ng:domain:project bundle-addons run
autostart=true
autorestart=true

numprocs=1
redirect_stderr=true
stdout_logfile=/home/app_billing/applications/purchase-gateway/pre-prod/logs/purchase-gateway-worker.log

[program:purchase-gateway-pre-prod-worker-legacy-import-rabbit-queue-one]
process_name=%(program_name)s_%(process_num)02d
command=php /home/app_billing/applications/purchase-gateway/pre-prod/current/artisan ng:message:consumer-purchase-to-legacy 1
autostart=true
autorestart=true

numprocs=1
# startsecs
# This is the minimum number of seconds that supervisor expects the process to run for it to be considered successful
# If the process stops in under the amount of seconds specified, supervisor will consider it exited with error and
# it will throw a FATAL.
startsecs=1
redirect_stderr=true
stdout_logfile=/home/app_billing/applications/purchase-gateway/pre-prod/logs/purchase-gateway-worker.log

[program:purchase-gateway-pre-prod-worker-member-profile-enriched-event-rabbit-queue-one]
process_name=%(program_name)s_%(process_num)02d
command=php /home/app_billing/applications/purchase-gateway/pre-prod/current/artisan ng:message:consumer-purchase-to-member-profile 1
autostart=true
autorestart=true

numprocs=1
# startsecs
# This is the minimum number of seconds that supervisor expects the process to run for it to be considered successful
# If the process stops in under the amount of seconds specified, supervisor will consider it exited with error and
# it will throw a FATAL.
startsecs=1
redirect_stderr=true
stdout_logfile=/home/app_billing/applications/purchase-gateway/pre-prod/logs/purchase-gateway-worker.log

[program:purchase-gateway-pre-prod-worker-business-group-sites]
process_name=%(program_name)s_%(process_num)02d
command=php /home/app_billing/applications/purchase-gateway/pre-prod/current/artisan  ng:domain:project business-group-sites run
autostart=true
autorestart=true

numprocs=1
redirect_stderr=true
stdout_logfile=/home/app_billing/applications/purchase-gateway/pre-prod/logs/purchase-gateway-worker.log

[program:purchase-gateway-pre-prod-worker-retry-failed-event-publishing]
process_name=%(program_name)s_%(process_num)02d
command=php /home/app_billing/applications/purchase-gateway/pre-prod/current/artisan ng:failed-event-publish:retry
autostart=true
autorestart=true

numprocs=1
redirect_stderr=true
stdout_logfile=/home/app_billing/applications/purchase-gateway/pre-prod/logs/purchase-gateway-worker.log

[program:purchase-gateway-pre-prod-worker-timer-pending-purchases]
process_name=%(program_name)s_%(process_num)02d
command=php /home/app_billing/applications/purchase-gateway/pre-prod/current/artisan ng:worker timer-pending-purchases run
autostart=true
autorestart=true

numprocs=1
# startsecs
# This is the minimum number of seconds that supervisor expects the process to run for it to be considered successful
# If the process stops in under the amount of seconds specified, supervisor will consider it exited with error and
# it will throw a FATAL.
startsecs=1
redirect_stderr=true
stdout_logfile=/home/app_billing/applications/purchase-gateway/pre-prod/logs/purchase-gateway-worker.log

[program:purchase-gateway-pre-prod-worker-void-transactions]
process_name=%(program_name)s_%(process_num)02d
command=php /home/app_billing/applications/purchase-gateway/pre-prod/current/artisan ng:worker void-transactions run
autostart=true
autorestart=true

numprocs=1
# startsecs
# This is the minimum number of seconds that supervisor expects the process to run for it to be considered successful
# If the process stops in under the amount of seconds specified, supervisor will consider it exited with error and
# it will throw a FATAL.
startsecs=1
redirect_stderr=true
stdout_logfile=/home/app_billing/applications/purchase-gateway/pre-prod/logs/purchase-gateway-worker.log
